# automation/launchd/com.oak.health-check.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.oak.health-check</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/robertnyborg/Projects/claude-oak-agents/automation/oak_notify.sh</string>
        <string>health</string>
    </array>

    <!-- Run every 3 days at 8am -->
    <key>StartInterval</key>
    <integer>259200</integer>  <!-- 3 days in seconds -->

    <key>StandardOutPath</key>
    <string>/Users/robertnyborg/Projects/claude-oak-agents/logs/health-check.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/robertnyborg/Projects/claude-oak-agents/logs/health-check.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>OAK_TELEMETRY_ENABLED</key>
        <string>true</string>
        <key>OAK_TELEMETRY_DIR</key>
        <string>/Users/robertnyborg/Projects/claude-oak-agents/telemetry</string>
        <key>PYTHONPATH</key>
        <string>/Users/robertnyborg/Projects/claude-oak-agents</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>LaunchOnlyOnce</key>
    <false/>
</dict>
</plist>
