# automation/launchd/com.oak.monthly-analysis.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.oak.monthly-analysis</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/robertnyborg/Projects/claude-oak-agents/automation/oak_notify.sh</string>
        <string>monthly</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Day</key>
        <integer>1</integer>  <!-- 1st of month -->
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/robertnyborg/Projects/claude-oak-agents/logs/monthly-analysis.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/robertnyborg/Projects/claude-oak-agents/logs/monthly-analysis.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>OAK_TELEMETRY_ENABLED</key>
        <string>true</string>
        <key>OAK_TELEMETRY_DIR</key>
        <string>/Users/robertnyborg/Projects/claude-oak-agents/telemetry</string>
        <key>PYTHONPATH</key>
        <string>/Users/robertnyborg/Projects/claude-oak-agents</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>LaunchOnlyOnce</key>
    <false/>
</dict>
</plist>
